# provisioning/muyu-invoice.service
[Unit]
Description=Muyu Invoice Generator
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=muyu
WorkingDirectory=/opt/muyu-invoice-generator
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=DATABASE_URL=postgres://invoice_user:invoice_pass@localhost:5432/invoice_db
ExecStart=/usr/bin/node src/web.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
